# ops/observer/observer.service
[Unit]
Description=MyApp off-site observer — DNS, TLS, full-path HTTP

[Service]
Type=oneshot
ExecStart=/etc/scripts/observer.sh
