# ops/myapp-static-gc.timer
[Unit]
Description=Weekly GC of unreferenced hashed static assets

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
