# ops/myapp-alert@.service
# Emails the operator that %i failed, via the relay in /etc/myapp/env.
# Instantiated by OnFailure= on the app, transactor, and watchdog units.
[Unit]
Description=Alert the operator: %i failed

[Service]
Type=oneshot
ExecStart=/etc/scripts/myapp-alert.sh %i
# --- Sandboxing.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
# ProtectSystem=strict makes /run read-only along with everything else, so
# the cooldown stamp myapp-alert.sh writes needs its slot GRANTED here —
# without it the alerter dies on the stamp write before it has sent
# anything, and the whole OnFailure= spine goes silent. Preserve= keeps
# the stamp dir (and so the half-hour cooldown) across oneshot runs; a
# reboot clears /run, which is the right reset — a reboot is news.
RuntimeDirectory=myapp-alert
RuntimeDirectoryPreserve=yes

