# ops/myapp.service
[Unit]
Description=MyApp — Clojure/Datomic SaaS
# The peer rides out a transactor blip on its own, but a FIRST boot on a
# fresh host needs the transactor listening before create-database! runs.
After=network-online.target myapp-transactor.service
Wants=network-online.target myapp-transactor.service
# A crash LOOP must converge to a message, not spin forever: Restart=
# heals blips, the start limit turns persistent failure into failed
# state, and failed state emails the operator.
OnFailure=myapp-alert@%n.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
User=myapp
Group=myapp
WorkingDirectory=/opt/myapp
Environment=MYAPP_PROFILE=prod
# systemd creates /var/log/myapp (owned by User=), and the security log
# (ch.42) lands there — the one file fail2ban tails and forensics greps.
LogsDirectory=myapp
Environment=MYAPP_SECURITY_LOG=/var/log/myapp/security.log
# Everything the :prod profile requires, in one root-owned 0600 file
# (see env.example). Boot refuses to start with any value missing.
EnvironmentFile=/etc/myapp/env
# Half the heap becomes the peer's object cache unless object-cache-max
# says otherwise — dev/bench.clj is where that cache's value was measured.
# ExitOnOutOfMemoryError turns heap exhaustion into a death Restart= can
# heal instead of a half-alive process nothing notices.
# The -D flag hands the peer's periodic metrics report to the app's
# /metrics endpoint (myapp.web.metrics) — object cache, memory, index
# activity. Wired here, not in code: a dev REPL never loads it.
ExecStart=/usr/bin/java -Xmx1g -XX:+ExitOnOutOfMemoryError -Ddatomic.metricsCallback=myapp.web.metrics/datomic-callback! -jar /opt/myapp/myapp.jar
Restart=on-failure
RestartSec=2
# stop delivers SIGTERM; the app drains in-flight requests for up to a
# second (its shutdown hook) — 15s of patience before systemd escalates.
TimeoutStopSec=15

# --- Sandboxing. The app talks to loopback and writes only its log dir;
# everything else is closed. NOTE: MemoryDenyWriteExecute is deliberately
# ABSENT — the JVM's JIT maps writable-then-executable pages and would
# die under it; that one directive is the JVM-shaped hole in the standard
# hardening set, and pretending otherwise would break the box.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
