# ops/myapp-transactor.service
[Unit]
Description=Datomic transactor — the single writer for myapp + analytics
After=network-online.target postgresql.service
Wants=network-online.target postgresql.service
OnFailure=myapp-alert@%n.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
User=datomic
Group=datomic
ExecStart=/opt/datomic/bin/transactor /etc/myapp/transactor.properties
Restart=on-failure
RestartSec=2

# Sandboxing. Also a JVM — MemoryDenyWriteExecute omitted (JIT).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
