# ops/myapp@.service
# Template unit for the minimal-downtime update pair: the instance name is
# the port. `systemctl start myapp@3001 myapp@3002` runs two builds side by
# side behind Caddy's upstream pool; the deploy script (deploy-pair.sh)
# starts the new instance, health-gates it, and drains the old one.
[Unit]
Description=MyApp — instance on port %i
After=network-online.target myapp-transactor.service
Wants=network-online.target myapp-transactor.service
# A crash LOOP must converge to a message, not spin forever: Restart=
# heals blips, the start limit turns persistent failure into failed
# state, and failed state emails the operator.
OnFailure=myapp-alert@%n.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
User=myapp
Group=myapp
WorkingDirectory=/opt/myapp
Environment=MYAPP_PROFILE=prod
Environment=MYAPP_PORT=%i
LogsDirectory=myapp
Environment=MYAPP_SECURITY_LOG=/var/log/myapp/security.log
EnvironmentFile=/etc/myapp/env
# One jar file PER INSTANCE, so the running instance's jar is never
# replaced underneath a live JVM — deploy-pair.sh copies the new build to
# the idle instance's path before starting it.
# Heap dump to the log dir on OOM (pid-named, so the two pair instances
# never collide) before the exit-and-restart erases the evidence. The
# steady-state unit (myapp.service) carries the loopback REPL; the pair is
# transient (only live during a deploy), so it does not — two instances
# cannot share one REPL port.
ExecStart=/usr/bin/java -Xmx1g -XX:+ExitOnOutOfMemoryError -XX:+HeapDumpOnOutOfMemoryError -XX:HeapDumpPath=/var/log/myapp/ -Ddatomic.metricsCallback=myapp.web.metrics/datomic-callback! -jar /opt/myapp/myapp-%i.jar
Restart=on-failure
RestartSec=2
TimeoutStopSec=15

# Sandboxing — identical to myapp.service; MemoryDenyWriteExecute omitted
# for the same JVM-JIT reason.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
