# ops/myapp-watchdog.service
[Unit]
Description=MyApp watchdog — health, mail relay, disk
OnFailure=myapp-alert@%n.service

[Service]
Type=oneshot
ExecStart=/etc/scripts/myapp-watchdog.sh
# --- Sandboxing.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true

